# packages/clip-lockdown/clip-lockdown/service/xccdf_review@.service
[Unit]
Description=XCCDF %i review
Before=systemd-user-sessions.service
ConditionFileIsExecutable=/bin/oscap

[Service]
User=root
Group=root
Type=simple
ExecStart=/bin/oscap xccdf eval --profile %i --results ${SSG_PATH}/${SYSTEM_NAME}-ssg-post-%i-results.xml --report ${SSG_PATH}/${SYSTEM_NAME}-ssg-post-%i-results.html --cpe ${CONTENT_PATH}/ssg-rhel7-cpe-dictionary.xml ${CONTENT_PATH}/${CONTENT_FILE}
RemainAfterExit=true
SuccessExitStatus=2

[Install]
WantedBy=multi-user.target
